# packages/daemon/launchd/app.pocket-t.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>app.pocket-t.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/pocket-t</string>
    <string>run</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/pocket-t-daemon.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/pocket-t-daemon.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>POCKET_T_RELAY_URL</key>
    <string>wss://relay.pocket-t.app</string>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
